# backup/etc/cron.d/remote-backup
0 */6 * * * ubackup /home/ubackup/backup-cron.sh
