# backup/etc/cron.d/remote-backup
0 */6 * * * root /home/ubackup/backup-cron-root.sh
